# ansible/oversearch/files/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
user=root
startsecs=5
startretries=5
stopsignal=TERM
stopwaitsecs=10
