# ansible/oversearch/files/etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
directory=/var/lib/redis
autostart=true
autorestart=true
user=redis
startsecs=5
startretries=5
stopsignal=TERM
stopwaitsecs=10
